# dot_config/exact_etc_mirror/systemd/system/nvidia-tdp.service
[Unit]
Description=Set NVIDIA power limit above default

[Service]
Type=oneshot
ExecStart=/usr/bin/nvidia-smi -pl 204
