# dot_config/etc_mirror/systemd/system/user@1000.service.d/delegate.conf
[Service]
Delegate=yes
# For ananicy, cgroups delegation
