# dot_config/exact_etc_mirror/udev/rules.d/99-rapl.rules
SUBSYSTEM=="powercap", KERNEL=="intel-rapl*", RUN+="/usr/bin/find /sys/class/powercap/intel-rapl/ -type f -exec chgrp -h kloud {} +", RUN+="/usr/bin/find /sys/class/powercap/intel-rapl/ -type f -exec chmod -h g+r {} +"
