# dot_config/etc_mirror/systemd/system/nvidia-tdp.timer
[Unit]
Description=Set power limit 10 seconds after boot

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target
